# root/etc/profile.d/03-alias.sh
alias ls='ls --color'
alias ll='ls -al'
alias aws='docker run -v ~/.aws:/root/.aws --rm -it amazon/aws-cli'
alias aws-whoami='aws sts get-caller-identity'
